# dunst/.config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst
NoDisplay=true
GenericName=Notification Daemon
Comment=Dunst Notification Daemon
Exec=dunst
